# systemd/.config/systemd/user/arch-update.service
[Unit]
Description=Notify about available Arch updates

[Service]
Type=oneshot
ExecStart=%h/.local/bin/arch-update-notify
