# systemd/.config/systemd/user/power-sound.service
[Unit]
Description=Power plug/unplug & low-battery sounds
After=graphical-session.target upower.service
Wants=graphical-session.target

[Service]
Type=simple
Environment=LOW_THRESH=15
ExecStart=%h/.local/bin/power-sound-monitor
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
