# systemd/.config/systemd/user/ppd-autoswitch.service
[Unit]
Description=Auto-switch power profile on AC/Battery via upower
After=graphical-session.target
Wants=graphical-session.target

[Service]
ExecStart=%h/.local/bin/ppd-autoswitch
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
