# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter for removable media
After=graphical-session.target
Wants=graphical-session.target

[Service]
ExecStart=/usr/bin/udiskie --tray --notify
Restart=on-failure

[Install]
WantedBy=default.target
